# deploy/systemd/sqr.service.example
[Unit]
Description=Sumbangan Query Rahmah
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/sumbanganqueryrahmah
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
TimeoutStopSec=30
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
TasksMax=256
MemoryMax=1G
ReadWritePaths=/var/www/sumbanganqueryrahmah/uploads /var/www/sumbanganqueryrahmah/output /var/www/sumbanganqueryrahmah/artifacts /var/www/sumbanganqueryrahmah/var

[Install]
WantedBy=multi-user.target
